# deploy/dispatch-api.service
# /etc/systemd/system/dispatch-api.service
[Unit]
Description=Field Service Dispatch Optimizer API (FastAPI + OR-Tools)
After=network.target postgresql.service

[Service]
User=www-data
WorkingDirectory=/opt/dispatch
Environment=DATABASE_URL=postgresql+psycopg://dispatch:CHANGEME@localhost:5432/dispatch
Environment=CORS_ORIGINS=https://dispatch.scottcampbell.io
# Seed (idempotent) before serving.
ExecStartPre=/opt/dispatch/.venv/bin/python data-generator/generate_synthetic_data.py
ExecStart=/opt/dispatch/.venv/bin/uvicorn backend.app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
